# roles/bigtop_hadoop/files/hadoop-yarn-timelineserver.service
[Unit]
Description=Hadoop YARN Timeline Server (ATS v2)
After=network.target
Wants=network.target

[Service]
Type=forking
User=hadoop
Group=hadoop
EnvironmentFile=/etc/default/hadoop
ExecStart=/usr/lib/hadoop-yarn/bin/yarn --daemon start timelinereader
ExecStop=/usr/lib/hadoop-yarn/bin/yarn --daemon stop timelinereader
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
